# systemd_services/alert_missed_blocks.service
[Unit]
Description=Qubetics missed block Telegram alert
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=admin
Group=admin
SupplementaryGroups=docker

# Use a directory that is guaranteed accessible
WorkingDirectory=/home/admin/scripts

# Provide env and fail early if it's missing
EnvironmentFile=/etc/default/telegram.env
ExecStartPre=/usr/bin/test -r /etc/default/telegram.env

# Call bash directly (avoid /usr/bin/env to keep logs simple)
ExecStart=/bin/bash /home/admin/scripts/alert_missed_blocks.sh

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Make sure no global hardening blocks /home (remove if you had these elsewhere)
# ProtectHome=false

[Install]
WantedBy=multi-user.target