# systemd_services/node-maintenance-tg-prestop.service
# /etc/systemd/system/node-maintenance-tg-prestop.service
[Unit]
Description=Send Telegram notice before node stops
DefaultDependencies=no
Before=shutdown.target halt.target reboot.target
Wants=network-online.target
After=network-online.target
ConditionPathExists=/home/admin/scripts/node_maintenance_tg.sh

[Service]
Type=oneshot
ExecStart=/home/admin/scripts/node_maintenance_tg.sh --pre-shutdown
User=root
Group=root
TimeoutStartSec=20

[Install]
WantedBy=halt.target reboot.target shutdown.target
