# systemd_services/validator-startup.service
# /etc/systemd/system/validator-startup.service
[Unit]
Description=Start validator node with Cosmovisor on boot
# Ensure Docker and binfmt are ready first
Requires=docker.service binfmt-amd64.service
After=network-online.target docker.service binfmt-amd64.service
Wants=network-online.target
# Make sure your data mount is available before running
RequiresMountsFor=/mnt/nvme/qubetics

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/admin/scripts/validator-startup.sh
WorkingDirectory=/home/admin
StandardOutput=journal
StandardError=journal
User=root

[Install]
WantedBy=multi-user.target
