# systemd_services/validator-startup.service
# /etc/systemd/system/validator-startup.service
[Unit]
Description=Start validator node (Docker or Cosmovisor) on boot
Wants=network-online.target docker.service
After=network-online.target docker.service
Requires=docker.service
RequiresMountsFor=/mnt/nvme/qubetics
ConditionPathExists=/home/admin/scripts/validator-startup.sh

[Service]
Type=oneshot
RemainAfterExit=yes
User=admin
WorkingDirectory=/home/admin
ExecStart=/home/admin/scripts/validator-startup.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target